# systemd/signage-player-lite.service
# Pi OS Lite kiosk service — pattern follows RPi Foundation's rpi-image-gen webkiosk example
# See .planning/KIOSK_COMPOSITOR_DECISION.md for rationale

[Unit]
Description=TinySignage Player (Kiosk Browser — Lite)
After=signage-app.service multi-user.target
Wants=signage-app.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=pi

# cage creates its own Wayland session on this TTY
TTYPath=/dev/tty1
StandardInput=tty-force
StandardOutput=journal
StandardError=journal

RuntimeDirectory=tinysignage
Environment=XDG_RUNTIME_DIR=/run/tinysignage

ExecStartPre=/bin/sleep 5
ExecStart=/usr/bin/cage -d -s -- /home/pi/TinySignage/venv/bin/python /home/pi/TinySignage/launcher.py
WatchdogSec=120
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
